# deployment/monitoring.service
[Unit]
Description=DGX Monitoring (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
User=accadmin
Group=shadow
SupplementaryGroups=shadow adm docker ollama devgroup
WorkingDirectory=/home/accadmin/workspace/monitoring
ExecStart=/usr/bin/python3 -m uvicorn backend.app.main:app --host 0.0.0.0 --port 3001
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
